# archlinuxcn/teleport-bin/teleport@.service
[Unit]
Description=Teleport Service
After=network.target

[Service]
Type=simple
Restart=on-failure
RestartSec=5
EnvironmentFile=-/etc/default/teleport-%i
ExecStart=/usr/bin/teleport start --config /etc/teleport-%i.yaml --pid-file=/run/teleport-%i.pid
ExecReload=/bin/sh -c "exec pkill -HUP -L -F /run/teleport-%i.pid"
PIDFile=/run/teleport-%i.pid
LimitNOFILE=524288

[Install]
WantedBy=multi-user.target
