# archlinuxcn/kernel-modules-hook/kernel-modules-hook.install
post_install() {
	cat << EOF

Please execute

  $ sudo systemctl daemon-reload
  $ sudo systemctl enable linux-modules-cleanup

Report any issues to:
https://github.com/saber-nyan/kernel-modules-hook/issues


										~desu

EOF
}
