# archlinuxcn/auto-cpufreq/auto-cpufreq.service
[Unit]
Description=auto-cpufreq - Automatic CPU speed & power optimizer for Linux
After=network.target network-online.target
Conflicts=power-profiles-daemon.service

[Service]
Type=simple
User=root
ExecStart=auto-cpufreq --daemon
Restart=on-failure

[Install]
WantedBy=multi-user.target
