# archlinuxcn/virtualbox-svn/vboxdrmclient.path
[Unit]
Description=VirtualBox Guest VMSVGA resize client
ConditionVirtualization=oracle

[Path]
PathExistsGlob=/dev/dri/renderD*
