# archlinuxcn/geph4-vpn-helper/geph4-vpn-helper.service
[Unit]
Description=Geph4 VPN Service
After=network.target nss-lookup.target
Conflicts=geph4-client.service

[Service]
EnvironmentFile=/etc/default/geph4-client
ExecStart=/usr/bin/geph4-vpn-helper /usr/bin/geph4-client $VPN_SUBCOMMAND $VPN_FLAGS $VPN_OPTIONS --credential-cache /tmp/geph4-credentials/ --password $VPN_PASSWORD --username $VPN_USERNAME
ExecStop=/usr/bin/geph4-vpn-helper
Restart=on-abort

[Install]
WantedBy=multi-user.target
