# archlinuxcn/keyd-gaomon-1061pro/keyd-gaomon-1061pro.install
post_install() {
    cat <<EOF
:: GAOMON 1061PRO Config Setup

1. Device Discovery:
   • Sandboxed: sudo keyd-bwrap monitor
   • Standard:  sudo keyd monitor

2. Enable Service (after config edit):
   • Sandboxed: sudo systemctl enable --now keyd-bwrap.service
   • Standard:  sudo systemctl enable --now keyd.service

Config path: /etc/keyd/gaomon-1061pro-tablet.conf
EOF
}

post_upgrade() {
    post_install
}
