# archlinuxcn/isso/isso.service
[Unit]
Description=Isso Comment Server
After=network.target network-online.target nss-lookup.target

[Service]
Environment="ISSO_SETTINGS=/etc/isso.conf"
ExecStart=/usr/bin/isso run

[Install]
WantedBy=multi-user.target
