# arch-audit/arch-audit.timer
[Unit]
Description=Run arch-audit once a day

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
