# archlinuxcn/nomad/nomad-server.service
[Unit]
Description=nomad server agent
Documentation=https://www.nomadproject.io/docs/agent/
After=network.target

[Service]
ExecStart=/usr/bin/nomad agent -config /etc/nomad/server.conf -server

[Install]
WantedBy=multi-user.target
