# archlinuxcn/daed/daed.install
post_install() {
	echo "---------------------------------------------------------------------------"
	echo "After installation, run \"sudo systemctl enable --now daed\" to start daed."
	echo "Then open your browser and navigate to http://localhost:2023"
	echo "---------------------------------------------------------------------------"
}
