# anbox-git/anbox-session-manager.service
[Unit]
Description=Anbox Session Manager

[Service]
ExecStart=/usr/bin/anbox session-manager

[Install]
WantedBy=default.target
