# anbox-git/anbox-container-manager.service
[Unit]
Description=Anbox Container Manager

[Service]
ExecStart=/usr/bin/anbox container-manager --daemon --privileged --data-path=/var/lib/anbox

[Install]
WantedBy=multi-user.target
