# archlinuxcn/himawaripy-git/himawaripy.install
post_install() {
	echo "To refresh your background every 10 minutes, run :"
	echo -e "\tsystemctl --user enable himawaripy.timer"
}
