# lantern/lantern.service
[Unit]
Description=Lantern
After=network.target

[Service]
ExecStart=/usr/bin/lantern -headless=true
User=lantern

[Install]
WantedBy=multi-user.target
