# minecraft-server/minecraftd-backup.timer
[Unit]
Description=Daily Minecraft Server Backup

[Timer]
OnCalendar=daily
AccuracySec=5min
Persistent=true

[Install]
WantedBy=multi-user.target
