# archlinuxcn/systemd-numlockontty/numlockontty.install
post_install() {
    cat << EOF
    ==> You should enable the service by running:
    # systemctl enable numLockOnTty

    To disable:
    # systemctl disable numLockOnTty
EOF
}

