# archlinuxcn/synergy/synergys.socket
[Unit]
Conflicts=synergys@.service

[Socket]
ListenStream=24800
Accept=false

[Install]
WantedBy=sockets.target
