# lxd/lxd.service
[Unit]
Description=REST API, command line tool and OpenStack integration plugin for LXC.
After=network.target
Requires=lxc-net.service

[Service]
Type=simple
ExecStart=/usr/bin/lxd --group lxd

[Install]
WantedBy=multi-user.target
