# archlinuxcn/cgproxy-git/cgproxy.install
#!/bin/sh

post_install(){
cat <<'DOC'
    to start service:
            systemctl enable --now cgproxy.service
DOC
}
