# archlinuxcn/openrazer/openrazer-daemon.install
msg() {
  echo "==> You might want to enable openrazer-daemon.service with the command \"systemctl --user enable --now openrazer-daemon.service\""
}

post_upgrade() {
  msg
}

post_install() {
  msg
}

