# archlinuxcn/graphite-api/graphite-api.socket
[Unit]
Description=Graphite API Socket

[Socket]
ListenStream=127.0.0.1:8888
ListenStream=[::1]:8888

[Install]
WantedBy=sockets.target
