# archlinuxcn/graphite-api/graphite-api.service
[Unit]
Description=Graphite-api daemon
After=network.target

[Service]
# listening address can be configured via a drop-in
Environment=LISTEN=127.0.0.1:8888 WORKERS=2
ExecStart=/usr/bin/gunicorn -w $WORKERS graphite_api.app:app -b $LISTEN
Type=simple

# some hardening
PrivateTmp=true
PrivateDevices=true
ProtectSystem=full
ProtectHome=true
NoNewPrivileges=true
CapabilityBoundingSet=
ReadOnlyPaths=/

[Install]
WantedBy=multi-user.target
