# archlinuxcn/caddy/caddy.install
post_install() {
	systemctl --quiet daemon-reload
	systemctl --quiet preset caddy.service
	mkdir -p /etc/caddy/caddy.conf.d
}

post_upgrade() {
	systemctl --quiet daemon-reload
	systemd-tmpfiles --create caddy.conf
	systemctl --quiet try-restart caddy.service
}

pre_remove() {
	systemctl --quiet --no-reload disable --now caddy.service
}

post_remove() {
	systemctl --quiet daemon-reload
}
