# qemu-git/qemu-ga.service
[Unit]
Description=QEMU Guest Agent

[Service]
ExecStart=/usr/bin/qemu-ga

[Install]
WantedBy=multi-user.target
