# archlinuxcn/adguardhome/adguardhome.service
[Unit]
Description=AdGuard Home: Network-level blocker
After=syslog.target network-online.target

[Service]
DynamicUser=true
StateDirectory=adguardhome
WorkingDirectory=/var/lib/adguardhome
AmbientCapabilities=CAP_NET_BIND_SERVICE CAP_NET_RAW
CapabilityBoundingSet=CAP_NET_BIND_SERVICE CAP_NET_RAW
ExecStart=/usr/bin/adguardhome -w /var/lib/adguardhome -l syslog

[Install]
WantedBy=multi-user.target
