# archlinuxcn/synergy-git/synergys_at.service
[Unit]
Description=Synergy Server Daemon
After=network.target

[Service]
User=%i
ExecStart=/usr/bin/synergys --no-daemon --config /etc/synergy.conf
Restart=on-failure

[Install]
WantedBy=multi-user.target
