# archlinuxcn/rait/rait@.service
[Unit]
Description=Redundant Array of Inexpensive Tunnels
After=network.target network-online.target nss-lookup.target

[Service]
Type=oneshot
ExecStart=/usr/bin/rait up --config=/etc/rait/%i.conf
ExecStop=/usr/bin/rait down --config=/etc/rait/%i.conf
ExecReload=/usr/bin/rait up --config=/etc/rait/%i.conf
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
