# archlinuxcn/rait/rait.service
[Unit]
Description=Redundant Array of Inexpensive Tunnels
After=network.target network-online.target nss-lookup.target

[Service]
Type=oneshot
ExecStart=/usr/bin/rait up
ExecStop=/usr/bin/rait down
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
