# archlinuxcn/kernel-modules-hook-reflink/kernel-modules-hook.install
post_install() {
	cat << EOF

Please execute

  $ sudo systemctl daemon-reload
  $ sudo systemctl enable linux-modules-cleanup

Report any issues to:
https://github.com/lideming/kernel-modules-hook/issues

This package is forked from:
https://github.com/saber-nyan/kernel-modules-hook


										~arigatou

EOF
}
