# archlinuxcn/seafile-server/seafile-server@.service
[Unit]
Description=Next-generation open source cloud storage with advanced features on privacy protection and teamwork.
After=syslog.target network.target

[Service]
Type=forking
WorkingDirectory=/srv/seafile/%i
ExecStart=sh -c 'seafile-server/seafile.sh start; seafile-server/seahub.sh start'
ExecStop=sh -c 'seafile-server/seafile.sh stop; seafile-server/seahub.sh stop'
PIDFile=/srv/seafile/%i/pids/seaf-server.pid
User=seafile
Group=seafile

[Install]
WantedBy=multi-user.target
