# archlinuxcn/vmware-workstation/vmware-networks.path
[Unit]
Description=Monitor to Load-On-Demand the VMware Networks service

[Path]
PathExistsGlob=/var/run/vmware/*/*

[Install]
WantedBy=paths.target
