# n2n-git/n2n-supernode.service
[Unit]
Description=N2N VPN Supernode
After=network.target

[Service]
Type=simple
User=nobody
ExecStart=/usr/share/n2n/supernode.sh

[Install]
WantedBy=multi-user.target

