# archlinuxcn/keyd-bwrap/keyd-bwrap.service
[Unit]
Description=key remapping daemon (bwrap sandboxed)

[Service]
Type=simple
ExecStart=/usr/bin/keyd-bwrap

[Install]
WantedBy=multi-user.target
