# rsync-noatime/rsyncd.socket
[Unit]
Description=Rsync Socket
Conflicts=rsyncd.service

[Socket]
ListenStream=873
Accept=true

[Install]
WantedBy=sockets.target
