# archlinuxcn/matrix-synapse-unstable/synapse.target
[Unit]
Description=Synapse parent target
After=network-online.target
Wants=network-online.target

[Install]
WantedBy=multi-user.target
