# archlinuxcn/rustdesk/rustdesk.install
# arg 1:  the new package version
post_install() {
  if systemctl -q is-enabled 'rustdesk.service'; then
    systemctl daemon-reload
    systemctl restart 'rustdesk.service'
  fi
}

# arg 1:  the new package version
# arg 2:  the old package version
post_upgrade() {
  if systemctl -q is-enabled 'rustdesk.service'; then
    systemctl restart 'rustdesk.service'
  fi
}
