# archlinuxcn/sunshine-roaldclark/sunshine.service
[Unit]
Description=Self-hosted game stream host for Moonlight
After=graphical.target
StartLimitIntervalSec=500
StartLimitBurst=5

[Service]
# Execution control
ExecStartPre=/usr/bin/sleep 3s
ExecStart=/usr/lib/sunshine/sunshine

# User and permissions
DynamicUser=yes
SupplementaryGroups=video input
CapabilityBoundingSet=CAP_SYS_ADMIN
AmbientCapabilities=CAP_SYS_ADMIN
NoNewPrivileges=yes

# Environment variables
Environment=XDG_CONFIG_HOME=/var/lib/sunshine
Environment=XDG_RUNTIME_DIR=/run/sunshine

# Directory settings
WorkingDirectory=/var/lib/sunshine
StateDirectory=sunshine
StateDirectoryMode=0700
RuntimeDirectory=sunshine
RuntimeDirectory=sunshine/pulse
RuntimeDirectoryMode=0700

# Path bindings
BindReadOnlyPaths=/run/user/1000/pulse/native:/run/sunshine/pulse/native

# Security sandboxing
ProtectSystem=strict
ProtectHome=yes
ProtectHostname=yes
ProtectClock=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
ProtectKernelLogs=yes
ProtectProc=invisible
ProcSubset=pid
PrivateTmp=yes
DevicePolicy=closed

# Device access
DeviceAllow=/dev/null r
DeviceAllow=/dev/urandom r
DeviceAllow=char-drm rw
DeviceAllow=char-input r
DeviceAllow=/dev/uinput rw
DeviceAllow=/dev/uhid rw
DeviceAllow=/dev/nvidia0 rw
DeviceAllow=/dev/nvidiactl rw
DeviceAllow=/dev/nvidia-modeset rw
DeviceAllow=/dev/nvidia-uvm rw

# Resource restrictions
RestrictRealtime=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
UMask=0077

# Restart policy
Restart=on-failure
RestartSec=5s
