# archlinuxcn/ebtables/ebtables.service
[Unit]
Description=Ethernet bridge table
Before=network-pre.target
Wants=network-pre.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/lib/systemd/scripts/ebtables start
ExecStop=/usr/lib/systemd/scripts/ebtables stop

[Install]
WantedBy=multi-user.target
