# archlinuxcn/nextcloud-systemd-timers/nextcloud-files-scan-all.service
[Unit]
Description=Nextcloud; scan all files

[Service]
User=nextcloud
ExecStart=/usr/bin/php -f /usr/share/webapps/nextcloud/occ files:scan --all

[Install]
WantedBy=basic.target
