# archlinuxcn/nextcloud-systemd-timers/nextcloud-files-scan-all.timer
[Unit]
Description=Nextcloud; scan all files weekly

[Timer]
OnBootSec=15min
OnCalendar=weekly
Persistent=True

[Install]
WantedBy=timers.target