# archlinuxcn/nextcloud-systemd-timers/nextcloud-cron.timer
[Unit]
Description=Nextcloud; run cron.php periodically

[Timer]
OnBootSec=5min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target