# archlinuxcn/nextcloud-systemd-timers/nextcloud-preview-generate-all.service
[Unit]
Description=Nextcloud; generate previews for all files

[Service]
User=nextcloud
ExecStart=/usr/bin/php -f /usr/share/webapps/nextcloud/occ preview:generate-all

[Install]
WantedBy=basic.target
