# archlinuxcn/vmware-horizon-client/vmware-horizon-virtual-printing.service
[Unit]
Description=VMware Horizon virtual printing
After=org.cups.cupsd.service

[Service]
Type=forking
Environment="THNUCLNT_SVC=global:daemon"
ExecStart=/usr/sbin/thnuclntd

[Install]
WantedBy=multi-user.target
