# alarmcn/planisphere-tools/planisphere-report.service
[Unit]
Description=Planisphere report

[Service]
Type=oneshot
ExecStart=/usr/bin/planisphere-report report
DynamicUser=yes
