# archlinuxcn/timeshift/timeshift.install
post_install() {
  if ! systemctl -q is-enabled 'cronie.service'; then
    echo 'TimeShift scheduled backups require the cronie.service to be running.'
    echo 'Try:'
    echo '  systemctl enable --now cronie.service'
  fi
}

post_upgrade() {
  post_install
}
