# archlinuxcn/pithos/systemd.service
[Unit]
Description=Pithos Media Player

[Service]
Type=dbus
BusName=io.github.Pithos
ExecStart=/usr/bin/pithos --gapplication-service
